# sys-process/thrash-protect/files/thrash-protect.rc8
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="thrash-protect daemon"
description=""
command=/usr/sbin/thrash-protect
pidfile=/var/run/thrash-protect.pid

start() {
  ebegin "Starting thrash-protect"
    start-stop-daemon --start --exec "${command}" --pidfile "${pidfile}" --make-pidfile --background --nicelevel -2
  eend $?
}
